# docker/config/etc/services.d/laravel-queue/run
#!/usr/bin/with-contenv bash

s6-setuidgid nginx ./.start-queue.sh
