# docker/config/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

/usr/sbin/crond -f -d 7
